# .config/systemd/user/osiris-manager.service
[Unit]
Description=OSiris manager WebSocket (8880)
After=network-online.target osiris-server.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/OSiris/osiris-service
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/OSiris/osiris-service/venv/bin/python -u /home/OSiris/osiris-service/osiris_manager.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
