# .config/systemd/user/speedvision-carscan.service
[Unit]
Description=CarScan realtime YOLO API (127.0.0.1:9001)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/OSiris/speedvision-service
Environment=PYTHONUNBUFFERED=1
Environment=CARSCAN_WEIGHTS=/home/OSiris/public_html/carscan/yolov8n.pt
# Optional: smaller input for faster CPU — Environment=CARSCAN_IMGSZ=416
ExecStart=/home/OSiris/speedvision-service/venv/bin/uvicorn server:app --host 127.0.0.1 --port 9001
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
